# dplaneos-v3.3.0/systemd/dplaned.service
[Unit]
Description=D-PlaneOS System Daemon
# Hard gate: dplaned MUST NOT start before ZFS pools are fully mounted.
# dplaneos-zfs-mount-wait.service polls until all pools are ONLINE + writable,
# preventing the Docker-before-ZFS silent data loss race condition.
After=network.target zfs.target dplaneos-zfs-mount-wait.service
Requires=dplaneos-zfs-mount-wait.service
Wants=zfs.target
Documentation=https://github.com/4nonX/homelab

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/run/dplaneos /var/lib/dplaneos /var/log/dplaneos /etc/dplaneos
ExecStartPre=/bin/chmod 755 /var/run/dplaneos
ExecStart=/opt/dplaneos/daemon/dplaned -db /var/lib/dplaneos/dplaneos.db -listen 127.0.0.1:9000 -smb-conf /var/lib/dplaneos/smb-shares.conf
WorkingDirectory=/opt/dplaneos
Restart=always
RestartSec=5

# Daemon runs as root — required for:
# - ZFS pool/dataset operations (zpool, zfs commands)
# - Disk power management (hdparm, smartctl)
# - Network configuration (ip, iptables)
# - Service management (systemctl reload smbd/nfsd)
# - udev socket communication (/var/run/dplaneos/)
User=root
Group=root

# RuntimeDirectory creates /run/dplaneos with correct ownership on boot
# and cleans it up on service stop — solves the udev socket permission issue
RuntimeDirectory=dplaneos
RuntimeDirectoryMode=0755
StateDirectory=dplaneos
LogsDirectory=dplaneos

# Hardening — restrict despite running as root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/dplaneos /var/lib/dplaneos /opt/dplaneos /etc/dplaneos /run/dplaneos /etc/crontab /etc/exports /etc/systemd/network
# Note: /etc/samba removed — NixOS owns smb.conf via modules/samba.nix
# Daemon writes per-share config to /var/lib/dplaneos/smb-shares.conf

# Capabilities (ZFS requires SYS_ADMIN, network config requires NET_ADMIN)
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_READ_SEARCH CAP_CHOWN CAP_FOWNER
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_DAC_READ_SEARCH CAP_CHOWN CAP_FOWNER

# Limits
LimitNOFILE=65536
TasksMax=4096

# OOM Protection — prevent daemon from killing the NAS
# 512MB should be more than enough for metadata management;
# if it exceeds this, something is wrong (memory leak, runaway indexing)
MemoryMax=512M
MemoryHigh=384M
OOMScoreAdjust=-900

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dplaned

[Install]
WantedBy=multi-user.target
