# install/systemd/dplaneos-init-db.service
[Unit]
Description=D-PlaneOS Database Initialization
Documentation=https://github.com/4nonX/D-PlaneOS
# Must run after local filesystems are mounted (DB lives on /var/lib)
After=local-fs.target
# Must complete before either daemon starts
Before=dplaned.service dplaneos-realtime.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root

# Ensure required directories exist with correct permissions
ExecStartPre=/bin/mkdir -p /var/lib/dplaneos /var/log/dplaneos /run/dplaneos
ExecStartPre=/bin/chmod 755 /var/lib/dplaneos /var/log/dplaneos /run/dplaneos

# Create DB file + enable WAL mode if not already present.
# Schema creation is handled by dplaned's initSchema() on first start -
# keeping schema in one place (schema.go) prevents drift.
ExecStart=/opt/dplaneos/install/scripts/init-database-with-lock.sh

# Validate indexes, WAL mode, and ANALYZE after init.
# Exits 0 even when warnings are found - non-fatal for first-run (tables created by daemon).
ExecStartPost=/opt/dplaneos/install/scripts/validate-db-schema.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dplaneos-init-db

[Install]
WantedBy=multi-user.target

