# install/systemd/dplaneos-realtime.service
[Unit]
Description=D-PlaneOS Real-time Event Daemon (stub — handled by dplaned)
Documentation=https://github.com/dplaneos/dplaneos
# Real-time monitoring and WebSocket are built into dplaned itself (/ws/monitor).
# This unit exists for compatibility; it becomes active immediately and exits cleanly.
After=dplaned.service
Wants=dplaned.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
